# data/ayatana-indicator-notifications.service.in
[Unit]
Description=Ayatana Indicator Notifications Service
PartOf=graphical-session.target
After=ayatana-indicators.target

[Service]
ExecStart=@pkglibexecdir@/ayatana-indicator-notifications-service
Restart=on-failure

[Install]
WantedBy=ayatana-indicators.target
